# deploy/systemd/net-monitor-cat.service
# 使用以下软链接方式快速创建 service 的引用与运行

# ln -s  $(pwd)/net-monitor-cat.service  /etc/systemd/system/net-monitor-cat.service
# 启动并设为开机自启
# sudo systemctl daemon-reload
# sudo systemctl enable net-monitor-cat.service
# sudo systemctl restart net-monitor-cat.service

[Unit]
Description=Network Monitor Service
# Documentation=file:/home/cat/uav/netmonitor/README.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/cat/uav/netmonitor
ExecStart=/home/cat/uav/netmonitor/net-monitor-arm
Restart=always
RestartSec=5

# Give the process a predictable place to create logs and read config from.
Environment=HOME=/home/cat/uav/netmonitor
Environment=NETMONITOR_CONFIG=/home/cat/uav/netmonitor/config/monitor-net.yaml

# Route switching requires network administration privileges.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
